# env_monitor.service
[Unit]
Description = twelite env monitor daemon

[Service]
ExecStart = /home/pi/TweliteEnvMonitor/env_monitor.py
Restart = always
Type = simple
User=pi
Group=pi

[Install]
WantedBy = multi-user.target
